# _github_3pack_safe/_staging/part3/docs/deploy/u25-scheduler.service.example
# /etc/systemd/system/u25-scheduler.service

[Unit]
Description=U25 auto scheduler (pipeline slots)
After=network.target u25-api.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/u25-scouting
EnvironmentFile=/opt/u25-scouting/.env
ExecStart=/opt/u25-scouting/.venv/bin/python u25_auto_scheduler.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
